# centos7/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/root/ssh.sh
